# erpnext16/single-aio/rootfs/etc/supervisor/conf.d/00-redis.conf
[program:redis]
command=/usr/local/bin/run-redis.sh
priority=20
autostart=false
autorestart=true
startretries=20
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis.err
